# cmd/crontab
# crontab in normal crontab format
#
# DO NOT EDIT THIS FILE
#
# create and edit the one in /spider/local_cmd/crontab
# for doing connections and things
#
2 0 * * * DXDebug::dbgclean()
2 30 * * 3 DXUser::export()
0 3 * * * Spot::daily()
